# lib/systemd/system/doorbelld.socket
[Unit]
Description=doorbell socket

[Socket]
ListenStream=0.0.0.0:1337
Accept=yes

[Install]
WantedBy=sockets.target
