# etc/systemd/dorfmap-bgdata.service
[Unit]
Description=Flukso Daemon
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/dorfmap-bgdata
User=www-data

[Install]
WantedBy=multi-user.target
