# etc/systemd/ring_bell.socket
[Unit]
Description=Ring the bell
After=network.target

[Socket]
ListenStream=1337
Accept=true

[Install]
WantedBy=multi-user.target
